# hai-tech-website.service
[Unit]
Description=Hai Tech Website Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/ameidar/.openclaw/workspace/projects/hai-tech-website
ExecStart=/home/ameidar/.nvm/versions/node/v22.22.0/bin/node server.js
Restart=always
RestartSec=3
Environment=NODE_ENV=production

[Install]
WantedBy=default.target
